# syncthing/supervisor_config.ini
[program:syncthing]
command = /usr/bin/syncthing -no-browser -home="/syncthing/config"
directory = /data
autorestart = True
user = root
environment = STNORESTART="1", HOME="/syncthing/config"