# syncthing/supervisor_config.ini
[program:syncthing]
command = /usr/bin/syncthing -no-browser -home="/home/pirate/.config/syncthing"
directory = /home/pirate/
autorestart = True
user = some_user
environment = STNORESTART="1", HOME="/home/pirate"