# syncthing/supervisor_config.ini
[program:syncthing]
command = /usr/bin/syncthing -no-browser -home="/home/syncthing/.config/syncthing"
directory = /home/pirate/
autorestart = True
user = syncthing
environment = STNORESTART="1", HOME="/home/syncthing"